# salt/master/salt-master_supervisor.conf
[program:salt-master]
command = salt-master
