# salt/minion/salt-minion_supervisor.conf
[program:salt-minion]
command = salt-minion
